# root/etc/services.d/firefox/run
#!/usr/bin/with-contenv bash

DISPLAY=${DISPLAY:-10}
export DISPLAY=":${DISPLAY}"

exec s6-setuidgid abc /usr/bin/firefox
